# schleuder.cron.weekly
#!/bin/sh

set -e
test -x /usr/bin/schleuder || exit 0

su -s /bin/sh schleuder -c "/usr/bin/schleuder check_keys"

exit 0
