# etc/schleuderd.service
[Unit]
Description=Schleuder API daemon

[Service]
ExecStart=/usr/local/bin/schleuderd
User=schleuder

[Install]
WantedBy=multi-user.target
