# debian/schleuder.cron.weekly
#!/bin/sh

test -x /usr/bin/schleuder || exit 0

su -s /bin/sh schleuder -c "/usr/bin/schleuder refresh_keys"
su -s /bin/sh schleuder -c "/usr/bin/schleuder check_keys"
